# openmct/packages/build/buildetc/new/etc/inittab
# main rc script
::sysinit:/etc/init.d/rcS

# Put a getty on the serial port
::respawn:/sbin/getty -L tts/0 115200 vt100 > /dev/null 2>&1

# Stuff to do before rebooting
null::shutdown:/bin/umount -a -r

# what to do when CTRL-ALT-DEL is pressed
#ca::ctrlaltdel:/sbin/shutdown -r -t 4 now
