# openmct/packages/build/buildetc/src/etc/inittab
::sysinit:/etc/rc.d/rcS
::respawn:/sbin/getty -L /dev/tts/0 115200 vt100
