# openmct/packages/build/buildetc/new/var/etc/rc.local
#!/bin/sh
/bin/echo "starte optionale Files..."
